# deployment/krishimitra-sensors.service
[Unit]
Description=KrishiMitra Sensor Hub (DHT22/BH1750/DS18B20/Arduino/LCD)
After=krishimitra-backend.service

[Service]
User=ess
WorkingDirectory=/home/ess/.openclaw/workspace/smart-agriculture-advisory-system
ExecStart=/usr/bin/python3 sensor_hub.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=BLINKA_FORCEBOARD=RASPBERRY_PI_4B

[Install]
WantedBy=multi-user.target
